# cron_files/mounts/mount-m.cron
#!/bin/bash

SCRIPT_NAME="mount-m.cron"

MOUNT_LOCATION="/mnt/data/movies"

MOUNT="dec-m:/"

MOUNT_CRYPT_SCRIPT_LOCATION="./mount-crypt"

. $MOUNT_CRYPT_SCRIPT_LOCATION