# cron_files/mounts/mount-tv.cron
#!/bin/bash

SCRIPT_NAME="mount-tv.cron"

MOUNT_LOCATION="/mnt/data/tv"

MOUNT="dec-tv:/"

MOUNT_CRYPT_SCRIPT_LOCATION="./mount-crypt"

. $MOUNT_CRYPT_SCRIPT_LOCATION